# redblueapp/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND

[program:php]
user=www-data
command=php -S 0.0.0.0:8888 -t /opt/net-diag